# supervisord.conf
[supervisord]
nodaemon=true

[program:spamassassin]
process_name = spamassassin
command = /etc/init.d/spamassassin start
startsecs = 0
autorestart = false

