# supervisord.conf
[supervisord]
nodaemon=true

[program:spampd]
process_name = spampd
command = /etc/init.d/spampd start
startsecs = 0
autorestart = false

